# deploy/reu-rating-monitor.service
[Unit]
Description=REU Rating Telegram Bot Hourly Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=reubot
Group=reubot
WorkingDirectory=/opt/reu-rating-bot/app
EnvironmentFile=/etc/reu-rating-bot.env
ExecStart=/usr/bin/env RUN_MODE=monitor_once /opt/reu-rating-bot/.venv/bin/python /opt/reu-rating-bot/app/bot.py
