# deploy/reu-rating-watchdog.service
[Unit]
Description=REU Rating Telegram Bot Watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=reubot
Group=reubot
EnvironmentFile=/etc/reu-rating-bot.env
ExecStart=/usr/bin/python3 /opt/reu-rating-bot/app/deploy/watchdog.py
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/reu-rating-bot/data
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
LockPersonality=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true

