# deploy/reu-rating-sync-retry.service
[Unit]
Description=REU Rating Telegram Bot Pending Sync Retry
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=reubot
Group=reubot
WorkingDirectory=/opt/reu-rating-bot/app
EnvironmentFile=/etc/reu-rating-bot.env
ExecStart=/usr/bin/env RUN_MODE=sync_pending_once /opt/reu-rating-bot/.venv/bin/python /opt/reu-rating-bot/app/bot.py
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/reu-rating-bot/data
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
LockPersonality=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true

